# Alpine/GoGS/roles/alpine_gogs/files/gogs.init.d
#!/sbin/openrc-run

USER=root
SERVICE="gogs"
COMMAND="/opt/gogs/src/github.com/gogits/gogs/gogs"
CMD_ARGS=${GOGS_OPTS}
DIR="/opt/gogs"

depend() {
    need net
    need localmount
}

start() {
    ebegin "Starting ${SERVICE}"
    start-stop-daemon --background --start \
    --make-pidfile --pidfile /run/${SERVICE}.pid \
    --exec ${COMMAND} \
    --user ${USER:-root} \
    --chdir ${DIR} \
    -- ${CMD_ARGS}
    eend $?
}

stop() {
    ebegin "Stopping ${SERVICE}"
    start-stop-daemon --stop --exec ${COMMAND} \
    --pidfile /run/${SERVICE}.pid
    eend $?
}