# Alpine/Jenkins/roles/alpine_jenkins/files/jenkins.init.d
#!/sbin/openrc-run

USER=jenkins
GROUP=jenkins
SERVICE="jenkins"
COMMAND="/usr/bin/java"
CMD_ARGS="${JENKINS_OPTS} -jar /opt/jenkins/jenkins.war"
DIR="/opt/jenkins"

depend() {
    need net
    need localmount
}

start() {
    ebegin "Starting ${SERVICE}"
    start-stop-daemon --background --start \
    --make-pidfile --pidfile /run/${SERVICE}.pid \
    --exec ${COMMAND} \
    --user ${USER:-root} \
    --group ${GROUP:-root} \
    --chdir ${DIR} \
    -- ${CMD_ARGS}
    eend $?
}

stop() {
    ebegin "Stopping ${SERVICE}"
    start-stop-daemon --stop --exec ${COMMAND} \
    --pidfile /run/${SERVICE}.pid
    eend $?
}